# src/crontab
0 8 * * * /src/script.sh >> /var/log/script.log